# supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
command=/bin/bash -c "service jenkins start"
